# etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command=/etc/init.d/elasticsearch start
stdout_logfile=/var/log/elasticsearch-stdout.log
redirect_stderr=true
